# systemd/ocw-health.timer
[Unit]
Description=Run OpenClaw Worker Health Check every 6h

[Timer]
OnCalendar=*-*-* 00,06,12,18:10:00
Persistent=true
RandomizedDelaySec=60
Unit=ocw-health.service

[Install]
WantedBy=timers.target
