# systemd/ocw-metrics.service
[Unit]
Description=OpenClaw Worker Metrics Snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/openclaw-worker/kit/bin/metrics_snapshot.sh
